# ChromiumOS/chromite/sdk/etc/profile.d/40-chromeos-cachedir.sh
export CHROMEOS_CACHEDIR="/var/cache/chromeos-cache"
